# conf/supervisor.conf
# -*- conf -*-
[program:neurogo]
command=/srv/sites/neurogo/conf/startup
user=ubuntu
group=ubuntu
autostart=true
autorestart=true
stderr_logfile=/srv/logs/neurogo/supervisor_err.log
